# scripts/rc/seulogin
#!/sbin/openrc-run

name="seulogin"
description="Southeast University Network Login Service"
command="/usr/bin/seulogin"
command_args="cron /etc/seulogin/${RC_SVCNAME#seulogin.}.toml"
pidfile="/run/${RC_SVCNAME}.pid"
command_background="yes"

depend() {
    need net
    after firewall
}

start_pre() {
    ebegin "Waiting for network to be fully online"
    sleep 60
    eend $?
}
